# dunst/services/battery.service
[Unit]
Description=Check whether the battery is running low

[Service]
ExecStart=/home/@USER/.local/bin/dunst/check-battery
