# dunst/services/boot.service
[Unit]
Description=Send notification after boot

[Service]
ExecStart=notify-send -t 5000 -u low -i boot "System booted" "Welcome aboard!"

[Install]
WantedBy=default.target