# dunst/services/temp.service
[Unit]
Description=Check if the CPU temperature is normal

[Service]
ExecStart=/home/@USER/.local/bin/dunst/check-temp
