# dunst/services/temp.timer
[Unit]
Description=Check if the CPU temperature is normal

[Timer]
OnBootSec=10min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
