# dunst/services/updates.service
[Unit]
Description=Check if there are updates available

[Service]
ExecStart=/home/@USER/.local/bin/dunst/check-updates
